# infra/supervisor/5estrelas-worker.conf
[program:5estrelas-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/5estrelas/app/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/5estrelas/app/storage/logs/worker.log
stopwaitsecs=3600

; UsuGer (lançador) → senior_cod_usu → departamento. Fila SOAP dedicada.
[program:5estrelas-worker-senior-launcher]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/5estrelas/app/artisan queue:work redis --queue=senior-launcher --sleep=1 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/www/5estrelas/app/storage/logs/worker-senior-launcher.log
stopwaitsecs=3600

; Fornecedor (cad_fornecedor). Fila SOAP separada — roda em paralelo com launcher.
[program:5estrelas-worker-senior-supplier]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/5estrelas/app/artisan queue:work redis --queue=senior-supplier --sleep=1 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/www/5estrelas/app/storage/logs/worker-senior-supplier.log
stopwaitsecs=3600

; Prontidão pós-enrich (só DB: dept + status). Leve — 2 workers.
[program:5estrelas-worker-senior-readiness]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/5estrelas/app/artisan queue:work redis --queue=senior-readiness --sleep=1 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/5estrelas/app/storage/logs/worker-senior-readiness.log
stopwaitsecs=3600

[program:5estrelas-reverb]
process_name=%(program_name)s
command=php /var/www/5estrelas/app/artisan reverb:start --host=0.0.0.0 --port=8080
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/5estrelas/app/storage/logs/reverb.log
stopwaitsecs=10

[program:5estrelas-schedule]
process_name=%(program_name)s
command=/bin/bash -c "while true; do php /var/www/5estrelas/app/artisan schedule:run --verbose --no-interaction >> /var/www/5estrelas/app/storage/logs/schedule.log 2>&1; sleep 60; done"
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/dev/null
